# mediawiki/cron/crontab
### MEDIAWIKI ###
# At every minute, run jobs
* * * * * /usr/local/bin/mwjobrunner
