# mediawiki/cron/crontab
### MEDIAWIKI ###
# At every minute, run jobs
* * * * * root /usr/local/bin/mwjobrunner
